# bin/routechoices-taskqueue.service
[Unit]
Description="routechoices django server"
Wants=network-online.target
After=network-online.target

[Service]
Type=Simple
ExecStart=/apps/routechoices-server/env/bin/python /apps/routechoices-server/manage.py process_tasks
Restart=always
Environment="PATH=/apps/routechoices-server/env/bin/"
WorkingDirectory=/apps/routechoices-server/
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutSec=5
PrivateTmp=true

[Install]
WantedBy=default.target
